# dornroeschen.service
[Unit]
Description=dornröschen (syncs NASen, backs up computers)

[Service]
ExecStart=/usr/bin/dornröschen \
	-sync=true \
	-backup=false
