# dornroeschen.timer
[Unit]
Description=Hourly Timer

[Timer]
OnCalendar=Mon-Fri, 10:00
Unit=dornroeschen.service

[Install]
WantedBy=multi-user.target
