# idlepoweroff.service
[Unit]
Description=power off once idle

[Service]
ExecStart=/usr/bin/idlepoweroff

[Install]
WantedBy=multi-user.target
